# usr/lib/lua/luci/view/themes/luci-theme-mytheme/root/etc/uci-defaults/luci-theme-iggitheme
#!/bin/sh
uci batch <<-EOF
	set luci.themes.iggitheme=/luci-static/iggitheme
	commit luci
EOF
